# openunum-gemini.service
[Unit]
Description=OpenUnum Gemini Version - Autonomous Assistant
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/corp-unum/OpenUnumGeminiVersion
ExecStart=/home/corp-unum/.bun/bin/bun run src/index.ts
Environment=PATH=/home/corp-unum/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=OPENUNUM_GEMINI_HOME=/home/corp-unum/.openunum-gemini
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
